# eval_server.supervisor.conf
[program:eval_server]
command=/root/workspace/.venv/bin/uvicorn eval_server:app --host 0.0.0.0 --port 9000
directory=/root/workspace
environment=TEUTONIC_EVAL_DATASET_MODE="raw_hippius",TEUTONIC_RAW_TOKENIZER_REPO="Qwen/Qwen3-4B",TEUTONIC_RAW_DATASET_MANIFEST="hf-mirrors/HuggingFaceFW/fineweb-edu/data/_manifest.json",TEUTONIC_RAW_DATASET_PREFIX="hf-mirrors/HuggingFaceFW/fineweb-edu/data",TEUTONIC_EVAL_N_PRIVATE="0",TEUTONIC_EVAL_N_PUBLIC="5000",EVAL_N_CAP="120000",EVAL_BATCH_SIZE="512",EVAL_MAX_RUNTIME_S="3600",TEUTONIC_TORCH_COMPILE="0",PYTHONUNBUFFERED="1"
autostart=true
autorestart=true
stderr_logfile=/var/log/eval_server.err.log
stdout_logfile=/var/log/eval_server.out.log
